# cloud/configs/celerybeat.conf
; Copyright (c) 2016-present, Facebook, Inc.
; All rights reserved.
;
; This source code is licensed under the BSD-style license found in the
; LICENSE file in the root directory of this source tree. An additional grant
; of patent rights can be found in the PATENTS file in the same directory.

[program:celerybeat]

directory=/var/www

command=envdir /var/opt/endagaweb-envdir celery beat -A endagaweb -l info -s /tmp/celerybeat-schedule --pidfile=/tmp/celerybeat.pid
user=ubuntu
numprocs=1
stdout_logfile=/var/log/celerybeat.log
stderr_logfile=/var/log/celerybeat.log
autostart=false
autorestart=true
startsecs=2

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 10

; When resorting to send SIGKILL to the program to terminate it
; send SIGKILL to its whole process group instead,
; taking care of its children as well.
killasgroup=true
