# cloud/configs/celeryd.conf
; Copyright (c) 2016-present, Facebook, Inc.
; All rights reserved.
;
; This source code is licensed under the BSD-style license found in the
; LICENSE file in the root directory of this source tree. An additional grant
; of patent rights can be found in the PATENTS file in the same directory.

[program:celery]

directory=/var/www

command=envdir /var/opt/endagaweb-envdir celery worker -A endagaweb -l info --autoreload --concurrency=10 -f /var/log/celery/%%n-%%i.log
user=ubuntu
numprocs=1
stdout_logfile=/var/log/worker.log
stderr_logfile=/var/log/worker.log
autostart=true
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 60

; When resorting to send SIGKILL to the program to terminate it
; send SIGKILL to its whole process group instead,
; taking care of its children as well.
killasgroup=true
